# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:avahi]
command=/usr/sbin/avahi-daemon --no-dbus -f /etc/avahi/avahi-daemon.conf
autostart=true
autorestart=true
priority=5
stopasgroup=true
killasgroup=true
stopwaitsecs=5
stdout_logfile=/var/log/supervisor/avahi.log
stderr_logfile=/var/log/supervisor/avahi.log

[program:samba-smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/smbd.log
stderr_logfile=/var/log/supervisor/smbd.log

[program:samba-nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/nmbd.log
stderr_logfile=/var/log/supervisor/nmbd.log

[program:nfs]
command=/usr/local/bin/run-nfs.sh
autostart=true
autorestart=true
priority=5
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stdout_logfile=/var/log/supervisor/nfs.log
stderr_logfile=/var/log/supervisor/nfs.log

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stopwaitsecs=5
stdout_logfile=/var/log/supervisor/vsftpd.log
stderr_logfile=/var/log/supervisor/vsftpd.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stopwaitsecs=5
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.log

[program:apache2]
command=/usr/sbin/apache2 -DFOREGROUND
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stopwaitsecs=10
environment=APACHE_RUN_USER=www-data,APACHE_RUN_GROUP=www-data,APACHE_RUN_DIR=/var/run/apache2,APACHE_LOG_DIR=/var/log/apache2,APACHE_PID_FILE=/var/run/apache2/apache2.pid
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2.log

[program:init-states]
command=/usr/local/bin/apply-services.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
priority=98
stdout_logfile=/var/log/supervisor/init-states.log
stderr_logfile=/var/log/supervisor/init-states.log

[program:webapp]
command=/app/venv/bin/gunicorn --bind 0.0.0.0:8080 --workers 2 --timeout 60 app:app
directory=/app
autostart=true
autorestart=true
priority=20
stopasgroup=true
killasgroup=true
stopwaitsecs=10
environment=CONFIG_DIR=/data/config
stdout_logfile=/var/log/supervisor/webapp.log
stderr_logfile=/var/log/supervisor/webapp.log
